# ftp_server/log_delete.service
[Unit]
Description=Delete Old Log Files

[Service]
Type=simple
ExecStart=/home/magma/ftp-server/log_delete.sh

